# deploy/stat.service
[Unit]
Description=Stat

[Service]
ExecStart=$HOME/stat/run.sh

[Install]
WantedBy=multi-user.target